# creative-intelligence/catscan-api.service
[Unit]
Description=Cat-Scan Creative Intelligence API
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/opt/catscan
Environment="PATH=/opt/catscan/venv/bin:/usr/bin"
Environment="CATSCAN_S3_BUCKET=rtbcat-csv-archive-frankfurt-328614522524"
Environment="CATSCAN_S3_REGION=eu-central-1"
Environment="CATSCAN_S3_ARCHIVE=true"
Environment="CATSCAN_RETENTION_DAYS=90"
ExecStart=/opt/catscan/venv/bin/python -m uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=catscan

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
